# samba/root/etc/services.d/samba-server/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid root /usr/sbin/smbd \
	--interactive --foreground -d 2 --configfile /config/smb.conf
