# squeezeboxserver/root/etc/services.d/lms-server/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid squeezeboxserver /usr/sbin/squeezeboxserver --prefsdir /config
